# scripts/systemd/wifi-telemetry.target
[Unit]
Description=WiFi Telemetry Service Group
Documentation=man:systemd.special(7)

[Install]
WantedBy=wifi-initialized.target
